# root/etc/services.d/weechatd/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc tmux new-session -d weechat -d /config
